# templates/default/init/alfresco.erb
#!/bin/sh
#
# chkconfig: 2345 80 30
# description: Alfresco One
#
### BEGIN INIT INFO
# Provides:          alfresco
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Should-Start:      $all
# Should-Stop:       $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start/stop Alfresco
# Description:       Start/stop Alfresco
### END INIT INFO

RETVAL=0

start () {
    su - <%= node['installer']['unix_user'] %> -c "/alfresco/4.2.0/alfresco.sh start \"$2\""
    RETVAL=$?
    if [ -d "/var/lock/subsys" ] && [ `id -u` = 0 ] && [ $RETVAL -eq 0 ] ; then
        touch /var/lock/subsys/alfresco
    fi

}

stop () {
    su - <%= node['installer']['unix_user'] %> -c "/alfresco/4.2.0/alfresco.sh stop \"$2\""
    RETVAL=$?
}


case "$1" in
    start)
        start "$@"
        ;;
    stop)
        stop "$@"
        ;;
    restart)
        stop "$@"
        start "$@"
        ;;
    *)
        /opt/alf-installation/alfresco.sh "$@"
        RETVAL=$?
esac
exit $RETVAL
